# rootfs/etc/inittab
::sysinit:/init
::restart:/sbin/init
tty1::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
